# bale-relay.service
[Unit]
Description=Bale Userbot Relay - Forward messages between Bale chats
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bale-relay
ExecStart=/opt/bale-relay/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
